# images/rpi-php7/etc/crontab
*/1 * * * * root echo my-example >> /var/log/cron.log 2>&1

